# autostart/wezterm.desktop
[Desktop Entry]
Type=Application
Name=WezTerm
Exec=wezterm start --cwd .
X-GNOME-Autostart-enabled=true
