# systemd/attendance-check.service
[Unit]
Description=bukuに登録した勤怠システムのURLを確認し、タブを開く/アクティブにする

[Service]
Type=oneshot
# 起動直後(Persistent=trueによるcatch-up実行)はセッションのPATHがまだ
# ~/.local/binを含んでいないことがあるため明示指定する
Environment=PATH=%h/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=%h/dotfiles/scripts/tab-check.sh attendance_check
