# systemd/notion-check.service
[Unit]
Description=bukuに登録したNotionのURLを確認し、タブを開く/アクティブにする

[Service]
Type=oneshot
ExecStart=%h/dotfiles/scripts/tab-check.sh notion_check
